# systemd-unit/ksysguard-gpu.service
[Unit]
Description=Collect GPU stats
After=network.target

[Service]
ExecStart=/usr/bin/ksysguard-gpu.py
Restart=on-failure

[Install]
WantedBy=default.target 
