# gpuStatsServer.service.default
[Unit]
Description=Collect GPU stats
After=network.target

[Service]
ExecStart=python3 /usr/bin/gpuStatsServer.py
Restart=on-failure

[Install]
WantedBy=default.target

